# launchd/io.hksw.codex-tap-sync.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.hksw.codex-tap-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>__REPO_ROOT__/scripts/run_update.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>CODEX_TAP_ENV_FILE</key>
    <string>__ENV_FILE__</string>
    <key>PYTHON3</key>
    <string>__PYTHON3__</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>StandardOutPath</key>
  <string>__LOG_FILE__</string>
  <key>StandardErrorPath</key>
  <string>__LOG_FILE__</string>
</dict>
</plist>
